# fs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
#ExecStart=-/sbin/agetty --autologin pi --noclear %I $TERM
ExecStart=
ExecStart=-/sbin/agetty --skip-login --noclear --noissue --login-options "-f pi" %I $TERM
